# deployment/trading-bot.service
[Unit]
Description=MVP Trading Bot - Crypto.com
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/MVP-NEW-TRADING-BOT
Environment="PATH=/home/ubuntu/MVP-NEW-TRADING-BOT/venv/bin"
ExecStart=/home/ubuntu/MVP-NEW-TRADING-BOT/venv/bin/python -m bot.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
